# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[program:xdummy]
command=/usr/bin/X :1 -configdir /etc/X11/xorg.conf.d -noreset -listen tcp
autostart=true
autorestart=true
priority=100
stderr_logfile=/var/log/supervisor/xdummy.err.log
stdout_logfile=/var/log/supervisor/xdummy.out.log

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=:1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox.err

[program:vnc]
command=/bin/bash -c 'for i in $(seq 1 30); do xdpyinfo -display :1 >/dev/null 2>&1 && break; sleep 0.5; done; exec x11vnc -display :1 -rfbport 10000 -rfbauth /root/.vnc/passwd -shared -forever -noxdamage -cursor most'
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc.err

[program:gemini-translator]
command=python main.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/gemini.log
stderr_logfile=/var/log/supervisor/gemini.err